# etc/bump.supervisor.conf
[program:bump]
command=sh run.sh
user=note
directory=/opt/bump
autorestart=true
